# raw/remote-mcp-launchd-live-proof/2026-04-20.com.hldpro.remote-mcp-monitor.rendered.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hldpro.remote-mcp-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/bennibarger/Developer/HLDPRO/hldpro-governance-remote-mcp-launchd-proof/scripts/remote-mcp/live_health_monitor.py</string>
    <string>--mode</string>
    <string>live</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/bennibarger/Developer/HLDPRO/hldpro-governance-remote-mcp-launchd-proof</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>StandardOutPath</key>
  <string>/Users/bennibarger/Developer/HLDPRO/hldpro-governance-remote-mcp-launchd-proof/projects/hldpro-governance/reports/remote-mcp-monitor.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bennibarger/Developer/HLDPRO/hldpro-governance-remote-mcp-launchd-proof/projects/hldpro-governance/reports/remote-mcp-monitor.stderr.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
