# launchd/com.hldpro.governance-observer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hldpro.governance-observer</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__REPO_ROOT__/scripts/orchestrator/read_only_observer.py</string>
    <string>--repo-root</string>
    <string>__REPO_ROOT__</string>
    <string>--report-root</string>
    <string>__REPO_ROOT__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/projects/hldpro-governance/reports/observer.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/projects/hldpro-governance/reports/observer.stderr.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
